# assets/service/log/run
#!/bin/sh
exec setuidgid vagrant multilog t ./main
